# shavak_etc/profile.d/examples.sh
export EXAMPLES=/usr/local/share/examples
