# shavak_etc/profile.d/sifdir.sh
export SIFDIR=/usr/local/share/singularity/images
